# deploy/systemd/metalgpt-web.service
[Unit]
Description=MetalGPT Web
# Запускается после vLLM и Docker (для Redis)
After=metalgpt-vllm.service docker.service
# Требует, чтобы vLLM и Docker были запущены
Requires=metalgpt-vllm.service docker.service

[Service]
User=metalgpt
WorkingDirectory=/opt/metalgpt/backend
EnvironmentFile=/opt/metalgpt/backend/.env
ExecStart=/opt/metalgpt/venv/bin/uvicorn main:app --host 127.0.0.1 --port 9000
Restart=always
# Добавляем задержку в 5 секунд перед перезапуском
RestartSec=5

[Install]
WantedBy=multi-user.target
