# deploy/systemd/metalgpt-web.service
[Unit]
Description=MetalGPT Web
After=metalgpt-vllm.service
Requires=metalgpt-vllm.service

[Service]
User=metalgpt
WorkingDirectory=/opt/metalgpt/backend
EnvironmentFile=/opt/metalgpt/backend/.env
ExecStart=/opt/metalgpt/venv/bin/uvicorn main:app --host 127.0.0.1 --port 9000
Restart=always
