# deploy/systemd/metalgpt-vllm.service
[Unit]
Description=MetalGPT vLLM
# Запускается после старта сети и Docker
After=network.target docker.service

[Service]
User=metalgpt
WorkingDirectory=/opt/metalgpt
EnvironmentFile=/opt/metalgpt/backend/.env
ExecStart=/opt/metalgpt/venv/bin/vllm serve nn-tech/MetalGPT-1 \
  --host 127.0.0.1 --port 8000 \
  --max-model-len 16384 \
  --gpu-memory-utilization 0.90 \
  --enable-chunked-prefill \
  --max-num-seqs 24 \
  --max-num-batched-tokens 12288 \
  --api-key ${VLLM_KEY}
Restart=always
# Добавляем задержку в 5 секунд перед перезапуском
RestartSec=5

[Install]
WantedBy=multi-user.target
